# app/files/tweaks/cmu-autorun/sdcard/recovery/02-run-tweaks-from-usb/99-run-tweaks.rules
ACTION=="add",KERNEL=="sd[abcdef]*", RUN+="/bin/sh /tmp/mnt/data_persist/dev/bin/02-run-tweaks-from-usb/udev_add_action_handler.sh %k"
